# configurations/app/app.ini
[program:porwhats]
directory=/app/web
command=bash -c 'python /app/app.py'
stdout_logfile=/app/logs/porwhats_supervisor.log
stderr_logfile=/app/logs/porwhats_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true