# configurations/app/app.ini
[program:porwhats]
directory=/app
command=bash -c 'python app.py'
stdout_logfile=/logs/porwhats_supervisor.log
stderr_logfile=/logs/porwhats_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true