# configurations/app/downloader.ini
[program:downloader]
directory=/python_downloader
command=bash -zc 'python downloader.py'
stdout_logfile=/logs/porwhats_downloader_supervisor.log
stderr_logfile=/logs/porwhats_downloader_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true