# configurations/app/downloader.ini
[program:downloader]
directory=/app/web
command=bash -c 'python downloader.py'
stdout_logfile=/app/logs/porwhats_downloader_supervisor.log
stderr_logfile=/app/logs/porwhats_downloader_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true