# configurations/app/whatsapp.ini
[program:whatsapp]
directory=/python_whats
command=bash -zc 'python whatsapp_stack_v1.py'
stdout_logfile=/logs/porwhats_whatsapp_supervisor.log
stderr_logfile=/logs/porwhats_whatsapp_supervisor_error.log
autostart=true
autorestart=true
startsecs=10
killasgroup=true
startretries=10
stopasgroup=true